# .config/systemd/user/battery-notify.service
[Unit]
Description=Notify If Battery Level is Low

[Service]
Type=oneshot
ExecStart=%h/.scripts/battery_notify
